# scripts/claude-reflections/launchd/com.claude.weekly-review.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.claude.weekly-review</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/travis.johnson/Development/BaptizeAmerica/LocalRepo/BaptizeWorld/scripts/claude-reflections/run-review.sh</string>
        <string>weekly</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>0</integer>  <!-- Sunday -->
        <key>Hour</key>
        <integer>20</integer>  <!-- 8 PM -->
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/travis.johnson/Development/BaptizeAmerica/LocalRepo/BaptizeWorld/scripts/claude-reflections/logs/weekly-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/travis.johnson/Development/BaptizeAmerica/LocalRepo/BaptizeWorld/scripts/claude-reflections/logs/weekly-stderr.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
